# overlay/systemd/antizapret.service.d/awg-reintegrate.conf
# AntiZapret-AWG: самовосстановление слоя AmneziaWG.
# Выполняется после КАЖДОГО старта antizapret.service (загрузка, авто-апдейт,
# ручной setup.sh). Drop-in лежит в /etc/systemd/system и переживает обновление
# AntiZapret (rm -rf /root/antizapret его не трогает). Префикс '-' — не фатально.
[Service]
ExecStartPost=-/opt/antizapret-awg/awg-reintegrate.sh
