# overlay/systemd/awg-quick@.service.d/override.conf
[Unit]
# awg-quick интерфейсы должны подняться ДО antizapret.service (маршрутизация)
Before=antizapret.service

[Service]
LimitNPROC=infinity
TasksMax=infinity
LimitNOFILE=infinity
