# bot/awg-stats.service
[Unit]
Description=AntiZapret-AWG stats collector (one poll)
After=awg-quick@antizapret.service

[Service]
Type=oneshot
Environment=AWG_STATS_DB=/root/antizapret/awg/stats.db
Environment=AWG_DIR=/etc/amnezia/amneziawg
ExecStart=/root/antizapret/awg/venv/bin/python /root/antizapret/awg/awg_stats.py poll
