# bot/awg-expire.service
[Unit]
Description=AntiZapret-AWG expire temporary clients

[Service]
Type=oneshot
ExecStart=/usr/local/bin/awg-client expire-check
