# bot/awg-stats.timer
[Unit]
Description=Poll AmneziaWG stats every 60s

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=10s

[Install]
WantedBy=timers.target
